# data/profile.d/make.sh
pkg: mkdir -p /etc/profile.d
pkg: cp -a %/terminal.sh /etc/profile.d
pkg: cp -a %/profile /etc
