# data/netatalk/com.saurik.afpd.dns-sd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.saurik.afpd.dns-sd</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/dns-sd</string>
        <string>-R</string>
        <string></string>
        <string>_afpovertcp._tcp</string>
        <string>local.</string>
        <string>548</string>
    </array>

    <key>OnDemand</key>
    <false/>
</dict>
</plist>
